# packaging/systemd/nova-metrics.service
[Unit]
Description=Nova Metrics Exporter (Prometheus)
Documentation=https://github.com/nova-project/nova
Requires=nova.service
After=nova.service

[Service]
Type=simple
ExecStart=/usr/bin/nova metrics serve --port 9090 --interval 15 --bind 0.0.0.0
Restart=on-failure
RestartSec=5s

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
